# .config/systemd/user/pcloud.service
[Unit]
Description=pCLoud

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=/home/hatchan/bin/pcloud

[Install]
WantedBy=default.target
